# deploy/hybrid-vpp-rl-supervisor.service
# User-level systemd template for the autonomous research supervisor.
# Copy to ~/.config/systemd/user/, adjust WorkingDirectory, then:
#   systemctl --user daemon-reload && systemctl --user enable --now hybrid-vpp-rl-supervisor
[Unit]
Description=hybrid-vpp-rl autonomous research supervisor
After=default.target

[Service]
Type=simple
WorkingDirectory=%h/hybrid-vpp-rl
ExecStart=%h/.local/bin/uv run python -m hybrid_vpp.training.supervisor
Restart=on-failure
RestartSec=30
Nice=5

[Install]
WantedBy=default.target
